# noncanon/implementation_bundle/systemd/actuator.service
[Unit]
Description=SLIME Dummy Actuator (NON-CANON)
After=network.target
Before=slime.service

[Service]
Type=simple
User=actuator
Group=slime-actuator
ExecStart=/usr/local/bin/actuator-dummy
Restart=no

# Runtime dir for the socket (optional, depending on your install script)
RuntimeDirectory=slime
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
