# noncanon/implementation_bundle/systemd/slime.service
[Unit]
Description=SLIME v0 (Law-Layer Actuation Gate)
After=network.target
Wants=actuator.service
After=actuator.service

[Service]
Type=simple

User=slime
Group=slime-actuator

ExecStart=/usr/local/bin/slime-runner

# Fail-closed: no automatic restart
Restart=no

# Runtime directory for /run/slime
RuntimeDirectory=slime
RuntimeDirectoryMode=0750

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal
SyslogIdentifier=slime

[Install]
WantedBy=multi-user.target

